# apps/proxy/resources/supervisord.conf
[supervisord]
nodaemon=true

[program:proxy]
redirect_stderr=true
stdout_logfile=/opt/output/index.html
command=/usr/sbin/nginx

[program:status]
command=/opt/run/simple_http_server.sh

